# conf/supervisor/supervisord.ini
[program:project]
command=uwsgi --emperor /path/to/project/conf/uwsgi/uwsgi_conf.ini
autostart=true
autorestart=true
stdout_logfile=/var/log/project/out.log
stderr_logfile=/var/log/project/err.log